# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:hadoop]
command=/usr/local/hadoop/sbin/start-dfs.sh
autostart=true
autorestart=false